# ops/openclaw/launchd/ai.openclaw.keepawake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.keepawake</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-ims</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/xvadur_mac/Library/Logs/openclaw/keepawake.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/xvadur_mac/Library/Logs/openclaw/keepawake.launchd.err.log</string>
</dict>
</plist>
